# uci-defaults/wireless

config wifi-device 'wlan0'
	option type 'mac80211'
	option channel '132'
	option country 'EU'
	option hwmode '11ah'

config wifi-iface 'radio0'
	option device 'wlan0'
	option network 'lan'
	option mode 'ap'
	option ssid 'EasyMesh'
	option encryption 'sae'
	option key '12345678'
	option hidden '0'
	option disabled '0'
	option multi_ap '3'

config wifi-device 'wlan1'
	option type 'mac80211'
	option channel '132'
	option country 'EU'
	option hwmode '11ah'

config wifi-iface 'radio1'
	option device 'wlan1'
	option network 'lan'
	option mode 'sta'
	option ssid 'EasyMesh'
	option encryption 'sae'
	option key '12345678'
	option hidden '0'
	option disabled '0'
	option multi_ap '3'
